# ops/systemd/hermes-hibernation-mcp-gateway.service
[Unit]
Description=Hermes Hibernation local MCP gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/aiwithapex/projects/hermes-lab
Environment=NODE_ENV=production
Environment=NODE_PATH=/home/aiwithapex/projects/hermes-lab/control-plane/node_modules
ExecStart=/home/aiwithapex/.nvm/versions/node/v24.14.0/bin/node /home/aiwithapex/projects/hermes-lab/ops/tunnel-client/hermes-mcp-gateway.cjs
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
